# resources/systemd/nemo_auto_cancel_training_sessions.timer
[Unit]
Description=Auto cancel training sessions with no registered users in NEMO

[Timer]
OnCalendar=minutely

[Install]
WantedBy=multi-user.target
